# systemd-services/airflow/airflow-worker.service
[Unit]
Description=Apache Airflow Worker
After=network.target

[Service]
User=<YOUR_USERNAME>
Group=<YOUR_GROUP>
Environment="AIRFLOW_HOME=/path/to/your/airflow"
ExecStart=/path/to/your/airflow/bin/airflow worker
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
